# deploy/systemd/rocom-egg-predictor.service
[Unit]
Description=Rocom Egg Predictor Next.js Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/rocom-egg-predictor
Environment=NODE_ENV=production
Environment=HOSTNAME=0.0.0.0
Environment=PORT=3000
Environment=NEXT_TELEMETRY_DISABLED=1
Environment=NODE_OPTIONS=--max-old-space-size=256
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
